# systemd/fcitx5.service
[Unit]
Description=Fcitx5

[Service]
ExecStart=fcitx5 --disable=wayland --verbose '*'=0
Restart=always

[Install]
WantedBy=default.target
